# config/systemd/fecoteme_access_control.service
[Unit]
Description=FECOTEME Hall Access Control System

[Service]
Type=simple
Restart=on-failure
RestartSec=3
ExecStart=/home/pi/fecoteme_access_control/code/fecoteme_access_control.py

[Install]
WantedBy=multi-user.target

